# supervisord/celerybeat.conf
; ================================
;  celery beat supervisor
; ================================

[program:celerybeat]
#### Config for FOR LOCAL Testing
command=/home/anarchistmegabyte/workspace/env/breakit-env/bin/celery  beat -A breakit:celery_app --schedule=/tmp/celerybeat-schedule --loglevel=INFO --pidfile=/tmp/celerybeat.pidcd
directory=/home/anarchistmegabyte/workspace/breakit

#### For Development Testing in Git branch Dev_Testing in Bitbucket
# command=/home/ubuntu/ebg-backend/bin/celery  beat -A ebgBackend:ebg_celery_app --schedule=/tmp/celerybeat-schedule --loglevel=INFO --pidfile=/tmp/celerybeat.pidcd
# directory=/home/ubuntu/ebikeGoBackend

numprocs=1
stdout_logfile=celerybeat.log
stderr_logfile=celerybeat.log
autostart=true
autorestart=true
startsecs=10
environment=VERSION='PROD'

; Causes supervisor to send the termination signal (SIGTERM) to the whole process group.
stopasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999
